# includes/AlarmPiHat.cron
# Run the AlarmPiHat poller every x minutes
*/5 * * * *     root  /usr/bin/php /home/pi/AlarmPiHat/poller.php
#
