# includes/AlarmPiHat.cron
# Run the AlarmPiHat poller every x minutes
*/1 * * * *  /usr/bin/php /home/pi/AlarmPiHat/poller.php
#
