# includes/AlarmPiHat.cron
# /etc/cron.d/AlarmPiHat: crontab fragment for AlarmPiHat
#  This runs the poller every x minutes.

# Run the poller every x minutes
*/5 * * * *     root  /home/pi/AlarmPiHat/poller.php
